# usr/local/share/legion-power-management/systemd/_common/power-ac.target
[Unit]
Description=On AC power
DefaultDependencies=no
Conflicts=power-battery.target
