# usr/local/share/legion-power-management/systemd/_common/power-battery.target
[Unit]
Description=On battery power
DefaultDependencies=no
Conflicts=power-ac.target
