# usr/local/share/legion-power-management/systemd/_common/power-profile-power-saver.target
[Unit]
Description=Power profile - Power-saver
DefaultDependencies=no
Conflicts=power-profile-performance.target power-profile-balanced.target
