# usr/local/share/legion-power-management/systemd/_common/power-profile-performance.target
[Unit]
Description=Power profile - Performance
DefaultDependencies=no
Conflicts=power-profile-power-saver.target power-profile-balanced.target
