# etc/systemd/system/power-battery.target
[Unit]
Description=On battery power
DefaultDependencies=no
